# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/app/run/supervisord.pid

[unix_http_server]
file=/app/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///app/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:nginx]
; Master runs as root so it can open /dev/stdout and /dev/stderr (owned by
; PID 1 root) for access/error logs. Worker processes drop to galactilog via
; the `user` directive in /etc/nginx/nginx.conf.
command=nginx -g "daemon off;"
priority=10
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:uvicorn]
command=uvicorn app.main:app --host 127.0.0.1 --port 8000 --proxy-headers
directory=/app
priority=20
autostart=true
autorestart=true
user=galactilog
environment=HOME="/tmp",USER="galactilog"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:celery]
command=celery -A app.worker.celery_app worker --loglevel=info --concurrency=%(ENV_GALACTILOG_CELERY_CONCURRENCY)s
directory=/app
priority=30
autostart=true
autorestart=true
user=galactilog
environment=HOME="/tmp",USER="galactilog"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:celery-beat]
command=celery -A app.worker.celery_app beat --loglevel=info --schedule /app/run/celerybeat-schedule
directory=/app
priority=25
autostart=true
autorestart=true
user=galactilog
environment=HOME="/tmp",USER="galactilog"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
